# etc/local.d/local.start
# This is a good place to load any misc programs
# on startup (use &>/dev/null to hide output)

kbdrate -r 50 -d 200 >/dev/null 2>&1
